# ROOT-COPY/etc/udev/rules.d/99-monitor-hotplug.rules
ACTION=="change", SUBSYSTEM=="drm", ENV{HOTPLUG}=="1", ENV{XAUTHORITY}="/home/empis/.Xauthority", RUN+="/home/empis/.scripts/monitorhotplug.sh"
